# factorio.service
[Service]
Type=forking
ExecStart=/opt/factorio/factorio_run.sh
PIDFile=/var/lib/factorio/factorio.pid
TimeStartSec=20
Restart=always